# user/chromium-over-mirai-socks-server@.service
[Service]
Environment=SOCKS_SERVER=localhost:49500
ExecStart=/usr/bin/chromium --user-data-dir=%f

Restart=always
