# user/chromium-over-proxy@.service
[Service]
Environment=SOCKS_SERVER=localhost:1080
ExecStart=/usr/bin/chromium --user-data-dir=%h/.config/chromium-over-proxy

Restart=always
