# user/chromium@.service
[Service]
ExecStart=/usr/bin/chromium --user-data-dir=%f

Restart=always
