# user/dwm.service
[Service]
ExecStart=/home/kalter/.local/bin/dwm
