# user/pavucontrol.service
[Service]
ExecStart=/usr/bin/pavucontrol
