# user/socks-server.service
[Service]
ExecStart=/usr/bin/ssh -N -D 1080 eu-north
Restart=on-failure
RestartSec=5s
