# user/st@.service
[Service]
ExecStart=/home/kalter/.local/bin/st
